# config/cron/rumi.cron
00 3 * * * cd /var/share/enisysgw; /usr/local/rvm/rubies/ruby-2.1.2/bin/ruby script/rails runner -e production 'Gw::Script::Schedule.delete'
10 3 * * * cd /var/share/enisysgw; /usr/local/rvm/rubies/ruby-2.1.2/bin/ruby script/rails runner -e production 'Gwbbs::Script::Task.delete'
20 3 * * * cd /var/share/enisysgw; /usr/local/rvm/rubies/ruby-2.1.2/bin/ruby script/rails runner -e production 'Gwcircular::Script::Task.delete'
30 3 * * * cd /var/share/enisysgw; /usr/local/rvm/rubies/ruby-2.1.2/bin/ruby script/rails runner -e production 'System::Script::Session.delete_expired_sessions'
40 3 * * * cd /var/share/enisysgw; /usr/local/rvm/rubies/ruby-2.1.2/bin/ruby script/rails runner -e production 'Gwbbs::Script::Task.preparation_delete'
50 3 * * * cd /var/share/enisysgw; /usr/local/rvm/rubies/ruby-2.1.2/bin/ruby script/rails runner -e production 'Doclibrary::Script::Task.preparation_delete'
00 4 * * * cd /var/share/enisysgw; /usr/local/rvm/rubies/ruby-2.1.2/bin/ruby script/rails runner -e production 'Gw::Script::Reminder.clear'
